# scripts/motd/update-motd.d/10-hostname.fish
figlet -f slant (hostname) | lolcat -tf
