# scripts/motd/update-motd.d/10-hostname.fish
figlet -f slant kingfisher | lolcat -tf
